# telegram-bridge.service.example
[Unit]
Description=Telegram Bridge — Claude CLI bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__BRIDGE_DIR__
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

# Env file with secrets
EnvironmentFile=__BRIDGE_DIR__/.env

# Graceful shutdown (bot handles SIGTERM)
KillSignal=SIGTERM
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-bridge

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=__HOME_DIR__ /tmp

[Install]
WantedBy=multi-user.target
